# etc/services.d/marktplaatsscraper/run
#!/usr/bin/with-contenv sh
exec python3 /config/Main.py -a ${PUSHOVER_API_TOKEN} -u ${PUSHOVER_USER} -i ${SCANNING_INTERVAL}
