# etc/services.d/marktplaatsscraper/run
#!/usr/bin/with-contenv sh
exec python3 /config/app/Main.py
